# lib/vtk/commands/socks/gov.va.socks.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key>
 <string><%= label %></string>
 <key>KeepAlive</key>
 <dict>
  <key>NetworkState</key>
  <true/>
 </dict>
 <key>RunAtLoad</key>
 <true/>
 <key>ProgramArguments</key>
 <array>
     <string><%= autossh_path %></string>
     <string>-v</string>
     <string>-M</string>
     <string>20001</string>
     <string>-D</string>
     <string><%= port %></string>
     <string>socks</string>
     <string>-N</string>
 </array>
 <key>StandardOutPath</key>
 <string><%= boot_script_path %>/Logs/autossh/autossh.stdout</string>
 <key>StandardErrorPath</key>
 <string><%= boot_script_path %>/Logs/autossh/autossh.stderr</string>
 <key>User</key>
 <string><%= user %></string>
 <key>ThrottleInterval</key>
 <integer>30</integer>
</dict>
</plist>
